# systemd/cf-owntracks.timer
[Unit]
Description=Periodic refresh of Cloudflare allowlist for OwnTracks
Documentation=file:///usr/local/share/cf-owntracks/README.md

[Timer]
# Every 6 hours with a randomized offset — ranges stay current with zero
# sysadmin involvement, without stampeding Cloudflare on the hour.
# The installer patches OnCalendar when --refresh-interval is given.
OnCalendar=*-*-* 00/6:00:00
RandomizedDelaySec=30m
# Catch up on missed runs (host was down/suspended).
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
