# systemd/cf-owntracks.service
[Unit]
Description=Refresh Cloudflare allowlist for OwnTracks (firewall + nginx)
Documentation=file:///usr/local/share/cf-owntracks/README.md
After=network-online.target nginx.service
Wants=network-online.target
# Transient failure? Try again in ~30 minutes instead of waiting for the
# next scheduled tick. The retry timer fires cf-owntracks.service once.
OnFailure=cf-owntracks-retry.timer

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cf-owntracks-refresh

# --- Hardening ---------------------------------------------------------------
# Root is required (firewall, nginx reload). The writable surface is narrowed
# via ProtectSystem=strict + an explicit ReadWritePaths allowlist.
# v2 fix: firewall persistence paths (/etc/nftables.d, /etc/nftables.conf,
# /etc/ufw, /etc/iptables) and the decision-log dir were missing in v1, which
# made persistence writes fail silently under the timer.
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
NoNewPrivileges=true
PrivateTmp=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
ReadWritePaths=/etc/nginx /etc/cf-owntracks /var/lib/cf-owntracks /etc/ssl/cloudflare /run /var/log /var/log/cf-owntracks /etc/nftables.d /etc/nftables.conf /etc/ufw /etc/iptables

TimeoutStartSec=300
Restart=no

[Install]
WantedBy=multi-user.target
